# systemd/user/base.target
[Unit]
Description=Required systemd user session services
Requires=dbus.socket
